# debian/postinst
#!/bin/sh
# postinst script for MythicOS

set -e

case "$1" in
    configure)

        # activer le service SysV si présent
        if command -v update-rc.d >/dev/null 2>&1; then
            update-rc.d mythicos-system defaults || true
        fi

        # compiler les schemas GLib
        if command -v glib-compile-schemas >/dev/null 2>&1; then
            glib-compile-schemas /usr/share/glib-2.0/schemas || true
        fi

        # recharger systemd et activer le service
        if command -v systemctl >/dev/null 2>&1; then
            systemctl daemon-reload || true
            systemctl enable mythicos-system.service || true
        fi

        # lancer l'ajustement MythicOS
        if [ -x /usr/lib/mythicos/mythicos-system/mythicos-adjust.py ]; then
            /usr/lib/mythicos/mythicos-system/mythicos-adjust.py || true
        fi

    ;;

    triggered)
        if [ -x /usr/lib/mythicos/mythicos-system/mythicos-adjust.py ]; then
            /usr/lib/mythicos/mythicos-system/mythicos-adjust.py || true
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

#DEBHELPER#

exit 0